# worker/regs_worker/regs-worker.service
[Unit]
Description=The Skinning Shed Regulations Worker
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/regs_worker
ExecStart=/usr/bin/node /opt/regs_worker/dist/index.js

# Graceful shutdown settings
TimeoutStopSec=180
KillSignal=SIGTERM
KillMode=mixed
SendSIGKILL=yes

# Restart settings
Restart=always
RestartSec=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=regs-worker

# Environment file
EnvironmentFile=/opt/regs_worker/.env

# Resource limits
LimitNOFILE=65535
MemoryMax=2G

[Install]
WantedBy=multi-user.target
